# scripts/setupcp/ao_templates/ubuntu16.04-target/home/.config/autostart/pidgin.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/pidgin
Hidden=false
X-MATE-Autostart-enabled=true
Name[en_US]=Pidgin
Name=Pidgin
Comment[en_US]=
Comment=
